# centos-django-uwsgi-supervisor-nginx/supervisor-app.ini
[program:app-uwsgi]
command = /usr/local/python3/bin/uwsgi --ini /home/docker/code/uwsgi.ini
startretries=1      ; 启动失败自动重试次数，默认是3

[program:nginx-app]
command = /usr/sbin/nginx